# wifi_config.service
[Unit]
Description=tH-Monitor WiFi Configuration Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/raspberry/tH-monitor
ExecStart=/usr/bin/python3 /home/raspberry/tH-monitor/wifi_safe_config.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
